# package/bob/cs-base-files/files/etc/crontabs/root
1 * * * * logrotate -v /etc/logrotate.d/rsyslog
